# deploy/workstation-agent.service
# systemd unit for the multi-user DCV on-box agent.
# Install: copy the binary to /usr/local/bin/workstation-agent and this unit to
# /etc/systemd/system/, then: systemctl enable --now workstation-agent
#
# The token verifier binds to loopback only and is reached over the client's SSM
# port-forward. Session provisioning (the /ensure-session endpoint, coming next)
# needs root to impersonate users for `dcv create-session`, hence User=root.

[Unit]
Description=Multi-user Amazon DCV workstation agent (session provisioning + token verifier)
After=dcvserver.service
Wants=dcvserver.service

[Service]
Type=simple
User=root
Environment=WSA_VERIFIER_ADDR=127.0.0.1:8444
ExecStart=/usr/local/bin/workstation-agent
Restart=on-failure
RestartSec=2
# Hardening: the agent needs root to manage sessions but nothing outside the host.
NoNewPrivileges=false
ProtectSystem=full

[Install]
WantedBy=multi-user.target
